# supervisord.conf
[supervisord]
user=root

[program:dockerd]
command=/usr/local/bin/dind dockerd --host=unix:///var/run/docker.sock --host=tcp://0.0.0.0:2375
redirect_stderr=true
stdout_logfile=/var/log/dind.log
