# rootfs/etc/services.d/dind/run
#!/usr/bin/execlineb -P

s6-notifyoncheck -c "docker version"
/usr/local/bin/dind dockerd
